# agent/magic-mirror-agent.service
[Unit]
Description=Magic Mirror PIR presence agent
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/magic-mirror/agent.env
ExecStart=/usr/local/bin/magic_mirror_agent.py
Restart=on-failure
RestartSec=3

User=magic-mirror-agent
Group=magic-mirror-agent
SupplementaryGroups=gpio

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
